# daemon/codex-usage-daemon.service
[Unit]
Description=Codexmeter usage daemon
After=bluetooth.target

[Service]
Type=simple
ExecStart=%h/codexmeter/daemon/codex-usage-daemon.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
